# systemd/kitty-theme-for-time-of-day.timer
[Unit]
Description=Run kitty-theme-for-time-of-day.service periodically

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
