# systemd/low-battery.service
[Unit]
Description=Send a notification if battery is running low.

[Service]
ExecStart=/home/teon/bin/notify-low-battery
