# systemd/low-battery.timer
[Unit]
Description=Run low-battery.service periodically

[Timer]
OnBootSec=300
OnUnitInactiveSec=300

[Install]
WantedBy=timers.target
