# systemd/sync-mail.timer
[Unit]
Description=Run sync-mail.service periodically

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
