# systemd/take-a-break.service
[Unit]
Description=Notify to take a break & drink water

[Service]
ExecStart=notify-send "Take a break & drink water"

